# deployment/supervisord/flux2.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
logfile=/var/log/supervisord.log
loglevel=info
pidfile=/var/run/supervisord.pid
nodaemon=false
childlogdir=/var/log/supervisor
environment=PATH="/opt/flux2/.venv/bin",PYTHONPATH="/opt/flux2/src"

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

# ============================================================================
# FLUX.2 UI WORKERS (Streamlit) - 3 Replicas with Round-Robin Load Balancing
# ============================================================================

[program:flux2-ui-worker-1]
command=/opt/flux2/.venv/bin/streamlit run ui_flux2_professional.py --server.port=8501 --server.address=127.0.0.1
directory=/opt/flux2
user=flux2
autostart=true
autorestart=true
startretries=3
startsecs=10
stopasgroup=true
priority=999
environment=CUDA_VISIBLE_DEVICES="0",PYTHONUNBUFFERED="1"
stdout_logfile=/var/log/supervisor/flux2-ui-worker-1.log
stderr_logfile=/var/log/supervisor/flux2-ui-worker-1-error.log
redirect_stderr=true

[program:flux2-ui-worker-2]
command=/opt/flux2/.venv/bin/streamlit run ui_flux2_professional.py --server.port=8502 --server.address=127.0.0.1
directory=/opt/flux2
user=flux2
autostart=true
autorestart=true
startretries=3
startsecs=10
stopasgroup=true
priority=999
environment=CUDA_VISIBLE_DEVICES="0",PYTHONUNBUFFERED="1"
stdout_logfile=/var/log/supervisor/flux2-ui-worker-2.log
stderr_logfile=/var/log/supervisor/flux2-ui-worker-2-error.log
redirect_stderr=true

[program:flux2-ui-worker-3]
command=/opt/flux2/.venv/bin/streamlit run ui_flux2_professional.py --server.port=8503 --server.address=127.0.0.1
directory=/opt/flux2
user=flux2
autostart=true
autorestart=true
startretries=3
startsecs=10
stopasgroup=true
priority=999
environment=CUDA_VISIBLE_DEVICES="0",PYTHONUNBUFFERED="1"
stdout_logfile=/var/log/supervisor/flux2-ui-worker-3.log
stderr_logfile=/var/log/supervisor/flux2-ui-worker-3-error.log
redirect_stderr=true

# ============================================================================
# MODEL WORKERS - Dedicated GPU Workers for Inference
# ============================================================================

[program:flux2-model-worker-4b]
command=/opt/flux2/.venv/bin/python scripts/model_worker.py
directory=/opt/flux2
user=flux2
autostart=true
autorestart=true
startretries=3
startsecs=10
stopasgroup=true
priority=998
environment=MODEL_WORKER_HOST="127.0.0.1",MODEL_WORKER_PORT="8600",MODEL_SERVICE_NAME="klein-4b-service",MODEL_KEY="flux.2-klein-4b",CUDA_VISIBLE_DEVICES="1",PYTHONUNBUFFERED="1"
stdout_logfile=/var/log/supervisor/flux2-model-worker-4b.log
stderr_logfile=/var/log/supervisor/flux2-model-worker-4b-error.log
redirect_stderr=true

[program:flux2-model-worker-9b]
command=/opt/flux2/.venv/bin/python scripts/model_worker.py
directory=/opt/flux2
user=flux2
autostart=true
autorestart=true
startretries=3
startsecs=10
stopasgroup=true
priority=998
environment=MODEL_WORKER_HOST="127.0.0.1",MODEL_WORKER_PORT="8601",MODEL_SERVICE_NAME="klein-9b-service",MODEL_KEY="flux.2-klein-9b",CUDA_VISIBLE_DEVICES="2",PYTHONUNBUFFERED="1"
stdout_logfile=/var/log/supervisor/flux2-model-worker-9b.log
stderr_logfile=/var/log/supervisor/flux2-model-worker-9b-error.log
redirect_stderr=true

# ============================================================================
# AUXILIARY SERVICES - Redis, RabbitMQ clients (if applicable)
# ============================================================================

[program:flux2-cache-manager]
command=/opt/flux2/.venv/bin/python -m flux2.redis_manager
directory=/opt/flux2
user=flux2
autostart=true
autorestart=true
startretries=2
startsecs=5
priority=997
environment=REDIS_HOST="127.0.0.1",REDIS_PORT="6379",PYTHONUNBUFFERED="1"
stdout_logfile=/var/log/supervisor/flux2-cache-manager.log
stderr_logfile=/var/log/supervisor/flux2-cache-manager-error.log
redirect_stderr=true

[program:flux2-queue-worker]
command=/opt/flux2/.venv/bin/python -m flux2.queue_worker
directory=/opt/flux2
user=flux2
autostart=true
autorestart=true
startretries=2
startsecs=5
priority=997
environment=RABBITMQ_HOST="127.0.0.1",RABBITMQ_PORT="5672",PYTHONUNBUFFERED="1"
stdout_logfile=/var/log/supervisor/flux2-queue-worker.log
stderr_logfile=/var/log/supervisor/flux2-queue-worker-error.log
redirect_stderr=true

# ============================================================================
# PROCESS GROUPS - For easier management
# ============================================================================

[group:flux2-ui]
programs=flux2-ui-worker-1,flux2-ui-worker-2,flux2-ui-worker-3
priority=100

[group:flux2-models]
programs=flux2-model-worker-4b,flux2-model-worker-9b
priority=99

[group:flux2-services]
programs=flux2-cache-manager,flux2-queue-worker
priority=98

# ============================================================================
# MONITORING & HEALTH CHECKS
# ============================================================================

[eventlistener:watchdog]
command=/opt/flux2/.venv/bin/python scripts/supervisor_watchdog.py
buffer_size=1000
events=PROCESS_STATE,TICK_60
redirect_stderr=true
stdout_logfile=/var/log/supervisor/watchdog.log

[inet_http_server]
port=127.0.0.1:9001

# ============================================================================
# NOTES
# ============================================================================
# 
# Installation:
#   pip install supervisor
#
# Start supervisord:
#   supervisord -c /opt/flux2/deployment/supervisord/flux2.conf
#
# Manage processes:
#   supervisorctl -c /opt/flux2/deployment/supervisord/flux2.conf status
#   supervisorctl -c /opt/flux2/deployment/supervisord/flux2.conf restart flux2:*
#   supervisorctl -c /opt/flux2/deployment/supervisord/flux2.conf stop flux2-ui-worker-1
#
# View logs:
#   tail -f /var/log/supervisor/flux2-ui-worker-1.log
#
# Enable auto-start on system boot (Linux):
#   sudo systemctl enable supervisord
#
